# systemd/estate-scraper-runpod-watchdog.timer
[Unit]
Description=Run estate-scraper RunPod pod auto-stop watchdog every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
